# etc/systemd/system/development-mode.service
[Unit]
Description=Sync skel with home
​
[Service]
ExecStart=/usr/bin/biglinux-sync-skel
