# etc/systemd/system/development-mode.path
[Unit]
Description="Monitor the /etc/skel file for changes"

[Path]
PathModified=/etc/skel/
Unit=development-mode.service

[Install]
WantedBy=multi-user.target
